# data/etc/systemd/system/lid-shutdown.service
[Unit]
Description=Shutdown on lid close

[Service]
Type=oneshot
ExecStart=/usr/local/bin/shutdown-lid.sh

[Install]
WantedBy=suspend.target